# supervisor-postgresql.conf
[supervisord]
nodaemon=true

[program:postgresql]
command = /usr/bin/run-postgresql.sh
user = postgres
autorestart = true
